# app/deb-build/etc/cron.d/pullbasedagent
0 0 * * * root /bin/bash /usr/local/pullbasedagent/src/cron_scripts/cleanup.sh
